# ds1821-update.timer
[Unit]
Description=Poll DS1821 temperature sensor every 60 seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
